# src/dhcpcd-online/dhcpcd-wait-online.service.in
[Unit]
Description=dhcpcd Wait Online
Documentation=man:dhcpcd-online(8)
Requisite=dhcpcd.service
After=dhcpcd.service
Before=network-online.target

[Service]
Type=oneshot
ExecStart=@BINDIR@/dhcpcd-online -q
RemainAfterExit=yes

[Install]
WantedBy=network-online.target
