# src/dhcpcd-online/dhcpcd-wait-online.service.in
[Unit]
Description=dhcpcd Wait Online
Requisite=dhcpcd.service
After=dhcpcd.service
Wants=network.target
Before=network.target network-online.target

[Service]
Type=oneshot
ExecStart=@BINDIR@/dhcpcd-online -q

[Install]
WantedBy=multi-user.target
